# packages/core/etc/systemd/system/simplevm-metadata-synchronizer.service
[Unit]
Description=Update metadata from metadata-server
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /etc/simplevm/utils/get_metadata.sh
WorkingDirectory=/etc/simplevm
EnvironmentFile=/etc/simplevm/metadata_config.env
User=root
StandardOutput=append:/var/log/metadata.log
StandardError=append:/var/log/metadata.log
Restart=on-failure

[Install]
WantedBy=multi-user.target