# packages/updater/DEBIAN/postinst
#!/bin/bash


# Enable the service and timer
systemctl enable simplevm-metadata-auto-update.service
systemctl enable simplevm-metadata-auto-update.timer

# Start the service and timer
systemctl start simplevm-metadata-synchronizer.timer

exit 0
