# packages/users/DEBIAN/postinst
!/bin/bash

set -e

# Reload systemd to recognize new service unit files
systemctl daemon-reload

# Enable the service and timer
systemctl enable public-keys-sync.service
systemctl enable public-keys-sync.timer

# Start the service and timer
systemctl start public-keys-sync.service
systemctl start public-keys-sync.timer

exit 0
